# crontab
28 10 * * * python -m robot --outputdir output/ birthdays.robot
