# crontab
28 10 * * * python -m robot birthdays.robot
